# server/starburger-clearsessions.timer
[Unit]
Description=Timer for starburger clear sessions

[Timer]
OnBootSec=300
OnUnitActiveSec=1w

[Install]
WantedBy=multi-user.target